# rootfs/etc/udev/rules.d/11-unbind_huawei_storage.rules
ACTION=="add", SUBSYSTEM=="scsi", DRIVERS=="sd", ATTRS{vendor}=="HUAWEI  ", RUN+="/home/pi/unbind_huawei_storage.py %k %p"
